# services.d/docker/run
#!/command/with-contenv sh
exec /usr/local/bin/dind dockerd --storage-driver=overlay2 $DOCKER_EXTRA_OPTS
