# services.d/docker/run
#!/usr/bin/with-contenv sh
set -e
exec /usr/local/bin/dind dockerd --storage-driver=overlay2 $DOCKER_EXTRA_OPTS
